# 45-disk.rules
SUBSYSTEM=="block", ACTION=="add", NAME="%k", RUN+="/root/diskadd.py"
