# 45-disk.rules
KERNEL=="sd[b-z][0-9]", SUBSYSTEM=="block", NAME="%k", RUN+="/root/diskmount/diskqueue.py", ENV{SYSTEMD_WANTS}="diskmount.service"
