# 45-disk.rules
KERNEL=="sd[b-z][0-9]", SUBSYSTEM=="block", NAME="%k", RUN+="/root/diskmount/diskadd.py"
