# diskmount.service
[Unit]
Description=Mounts connected drives

[Service]
ExecStart=/root/diskmount/diskupdate.py
